# rt.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/share/plcrt
HOME=/usr/share/plcrt

*/5 * * * * root rtcron.sh

